# src/dockerfiles/zeek/supervisord.conf
[supervisord]
nodaemon=true
[program:zeek]
command=/opt/zeek/bin/zeekctl deploy
